# systemd_services/godns.service
[Service]
Type=simple
ExecStart=/opt/godns/godns
User=godns
Group=godns

[Install]
WantedBy=multi-user.target