# systemd_services/imgwatch.service
[Service]
Type=simple
ExecStart=/opt/imgwatch.sh
ExecReload=/bin/kill -HUP $MAINPID
WorkingDirectory=/opt

[Install]
WantedBy=multi-user.target
